# packaging/debian/postinst
#!/bin/sh
set -e

if ! getent group docker > /dev/null 2>&1; then
    groupadd --system docker || true
fi

if ! id -u rustploy > /dev/null 2>&1; then
    useradd --system \
        --no-create-home \
        --shell /usr/sbin/nologin \
        --groups docker \
        --comment "Rustploy daemon" \
        rustploy
fi

# Ensure /etc/rustploy/ is owned by rustploy
chown -R rustploy:rustploy /etc/rustploy 2>/dev/null || true
chmod 750 /etc/rustploy 2>/dev/null || true

# Generate the RWP remote-admin block once, with a random token but disabled
# by default. Only on a genuine first install (dpkg passes an empty $2 when
# there's no previously-configured version) — NOT merely "section absent",
# so an admin who deletes [rwp] after install never gets it silently
# reappended on the next upgrade/reinstall.
CONFIG=/etc/rustploy/config.toml
if [ -z "$2" ] && [ -f "$CONFIG" ] && ! grep -q '^\[rwp\]' "$CONFIG"; then
    TOKEN=$(LC_ALL=C tr -dc 'A-Za-z0-9' < /dev/urandom | head -c 16)
    cat >> "$CONFIG" <<EOF

# Rustploy Wire Protocol — remote administrative channel (TCP).
# Disabled by default. To allow remote clients: set enabled = true, and to
# expose beyond loopback set bind_address = "0.0.0.0" (a token is required
# for any non-loopback bind). Token below was generated at install time.
[rwp]
enabled = false
bind_address = "127.0.0.1"
port = 8787
token = "$TOKEN"
max_connections = 8        # opcional (default 8)
max_frame_size = 1048576   # opcional (default 1 MiB)
read_timeout_secs = 15     # opcional
idle_timeout_secs = 120    # opcional
EOF
fi

# Append [env_backup] block on first install (when section is absent).
# Snapshots are written to <db_path>/env_backups/ every 60 s by default.
# See docs/env-backup.md for details and restore instructions.
if [ -f "$CONFIG" ] && ! grep -q '^\[env_backup\]' "$CONFIG"; then
    cat >> "$CONFIG" <<'EOF'

# Env var backup — snapshots automáticos de todas as env vars de projetos e serviços.
# Restaurar: rustploy env-backup list | rustploy env-backup restore <snapshot>
# Docs: /usr/share/doc/rustploy/env-backup.md
[env_backup]
# dir = "/var/lib/rustploy/db/env_backups"  # default: <db_path>/env_backups/
interval_secs = 60
EOF
fi

# config.toml now holds the RWP token — keep it readable only by rustploy.
chown rustploy:rustploy "$CONFIG" 2>/dev/null || true
chmod 640 "$CONFIG" 2>/dev/null || true

if [ "$1" = "configure" ]; then
    systemctl daemon-reload              >/dev/null 2>&1 || true
    systemctl enable  rustployd          >/dev/null 2>&1 || true
    # Helper privilegiado de firewall (socket activation): o daemon pede
    # allow/deny de portas externas via /run/rustploy/fw.sock.
    systemctl enable  rustployd-fw.socket >/dev/null 2>&1 || true
    systemctl restart rustployd-fw.socket >/dev/null 2>&1 || true
    systemctl restart rustployd          >/dev/null 2>&1 || true
fi

#DEBHELPER#
