# packaging/rustployd.service
[Unit]
Description=Rustploy PaaS Daemon
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=rustploy
Group=rustploy
SupplementaryGroups=docker
ExecStart=/usr/bin/rustployd
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

RuntimeDirectory=rustploy
RuntimeDirectoryMode=0755
# Sem isso, `systemctl restart rustployd` apaga /run/rustploy inteiro no stop
# (comportamento padrão de RuntimeDirectory=) e leva junto o fw.sock que
# rustployd-fw.socket já tinha vinculado ali — o socket some do disco mas
# `systemctl status rustployd-fw.socket` continua reportando active/listening
# (o fd do systemd segue aberto, só a entrada no diretório é que sume).
RuntimeDirectoryPreserve=yes
StateDirectory=rustploy
ConfigurationDirectory=rustploy

Environment=RUSTPLOY_SOCKET_PATH=/run/rustploy/rustploy.sock
Environment=RUSTPLOY_DB_PATH=/var/lib/rustploy/db
Environment=RUST_LOG=info
EnvironmentFile=-/etc/rustploy/environment

StandardOutput=journal
StandardError=journal
SyslogIdentifier=rustployd

NoNewPrivileges=yes
PrivateTmp=yes
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
